# apache-guacamole/test/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x720x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -rfbport 5900 -nopw -forever -shared
autorestart=true

[program:hello]
command=/usr/bin/python3 /opt/gui/hello_gui.py
autorestart=true

